# dpkg/DEBIAN/prerm
#!/bin/bash

systemctl disable --now radxa-usbnet.service || true